# containers/nginx-results/crontab
* * * * * /scripty/clean_results.py
# Mandatory blank line